# iso/live-build-tmp/config/includes.chroot/etc/systemd/system/lucid-console.service
[Unit]
Description=LUCID EMPIRE TITAN Web Console Dashboard
Documentation=https://github.com/ranpatidewage6-maker/lucid-empire-linux
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
Environment="TITAN_HOME=/opt/lucid-empire"
Environment="FLASK_ENV=production"
WorkingDirectory=/opt/lucid-empire/console
ExecStart=/usr/bin/python3 /opt/lucid-empire/console/app.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=false
ProtectSystem=false
PrivateTmp=true

[Install]
WantedBy=multi-user.target
