# iso/live-build-tmp/config/includes.chroot/etc/systemd/system/lucid-titan.service
[Unit]
Description=Lucid Titan V5 Identity Daemon
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/lucid-empire
ExecStart=/usr/bin/python3 /opt/lucid-empire/titan_core.py
Restart=always
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
